# systemd/sc-drop-watcher.service
[Unit]
Description=sc-drop-watcher — RSI store / comm-link drop notifier
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/projects/sc-drop-watcher
EnvironmentFile=/home/ubuntu/projects/sc-drop-watcher/.env
ExecStartPre=/bin/sh -c '[ -x /home/ubuntu/.nvm/versions/node/v24.13.1/bin/node ] || { echo "Pinned node binary not found at /home/ubuntu/.nvm/versions/node/v24.13.1/bin/node — update ExecStart in systemd/sc-drop-watcher.service after the nvm Node upgrade, then daemon-reload"; exit 1; }'
# Use `node --import tsx` instead of the `tsx` CLI. The CLI spawns a helper
# subprocess and uses an AF_UNIX IPC pipe (/tmp/tsx-<uid>/<pid>.pipe) to talk
# to it; the RestrictAddressFamilies= directive below blocks AF_UNIX and the
# CLI crashes with EAFNOSUPPORT. The --import form loads tsx as an ESM hook
# in-process — no IPC, no subprocess, no AF_UNIX needed.
ExecStart=/home/ubuntu/.nvm/versions/node/v24.13.1/bin/node --import tsx --env-file-if-exists=.env /home/ubuntu/projects/sc-drop-watcher/src/main.ts
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

# Hardening — this is a read-only HTTP client, no need for FS write outside its dir.
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/ubuntu/projects/sc-drop-watcher
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_INET AF_INET6
LockPersonality=true
MemoryDenyWriteExecute=false
RestrictRealtime=true

[Install]
WantedBy=multi-user.target
